# social/base/jaeger/jaeger
#!/sbin/openrc-run

name="jaeger"
command="/go/bin/all-in-one-linux"
command_args="--sampling.strategies-file /etc/jaeger/sampling_strategies.json"
command_background=true
command_user="jaeger"
pidfile="/run/rc_${RC_SVCNAME}.pid"

depend() {
  after *
}
